# etc/supervisor/conf.d/rsyslogd.conf

[program:rsyslog]
priority=0
command = rsyslog -n
autostart=true
autorestart=true
